# etc/systemd/system/minecraft.service
[Unit]
Description=Minecraft Server (Optimized for Pi 4 - 8GB)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=minecraft
Group=minecraft
WorkingDirectory=/home/minecraft/minecraft-server

# Optimized Java startup with Aikar's flags for 3GB RAM
ExecStart=/usr/bin/java -Xms3G -Xmx3G \
  -XX:+AlwaysPreTouch \
  -XX:+DisableExplicitGC \
  -XX:+ParallelRefProcEnabled \
  -XX:+PerfDisableSharedMem \
  -XX:+UnlockExperimentalVMOptions \
  -XX:+UseG1GC \
  -XX:G1HeapRegionSize=8M \
  -XX:G1HeapWastePercent=5 \
  -XX:G1MaxNewSizePercent=40 \
  -XX:G1MixedGCCountTarget=4 \
  -XX:G1MixedGCLiveThresholdPercent=90 \
  -XX:G1NewSizePercent=30 \
  -XX:G1ReservePercent=20 \
  -XX:InitiatingHeapOccupancyPercent=15 \
  -XX:MaxGCPauseMillis=200 \
  -XX:MaxTenuringThreshold=1 \
  -XX:SurvivorRatio=32 \
  -Dusing.aikars.flags=https://mcflags.emc.gs \
  -Daikars.new.flags=true \
  -jar paper-1.21.10-129.jar nogui

# Standard input/output
StandardInput=null
StandardOutput=journal
StandardError=journal

# Graceful shutdown - sends stop command to server
SuccessExitStatus=0 1
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=120

# Restart on failure
Restart=on-failure
RestartSec=15

# Resource limits (safety margins)
CPUQuota=250%           # Max 2.5 CPU cores
MemoryMax=3.5G          # Hard limit 3.5GB RAM
LimitNOFILE=4096        # File descriptor limit

[Install]
WantedBy=multi-user.target
